# home/autostart/yo-bitch.desktop
[Desktop Entry]
Type=Application
Exec=yo-bitch
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name=YoBitch
